# deploy/eve-memory-daily.timer
[Unit]
Description=Eve memory rollup (daily) — ~04:00 local

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
