# deploy/eve-memory-weekly.timer
[Unit]
Description=Eve memory rollup (weekly) — Sun ~04:15 local

[Timer]
OnCalendar=Sun *-*-* 04:15:00
Persistent=true

[Install]
WantedBy=timers.target
