# deploy/eve-memory-yearly.timer
[Unit]
Description=Eve memory rollup (yearly) — Jan 1 ~04:25 local

[Timer]
OnCalendar=*-01-01 04:25:00
Persistent=true

[Install]
WantedBy=timers.target
