# deploy/iva-telegram-poll.service
# systemd USER unit — Telegram polling-мост (бот отвечает без webhook/прокси).
#
# Сам забирает апдейты у Telegram (getUpdates) и POST-ит их в локальный eve.
# Плейсхолдеры __PROJECT_DIR__/__NODE_BIN__ подставляет install.sh.
# Установка вручную: см. комментарии в deploy/iva-memory-daily.service.
[Unit]
Description=Iva Telegram polling bridge (getUpdates → local webhook)
After=network-online.target iva.service
Wants=iva.service

[Service]
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__NODE_BIN__ --env-file=.env scripts/telegram-poll.mjs
# Bootstrap units introduced by an update: the updater that fetched this file may still
# be running the previous CLI code and therefore cannot know about the new timer yet.
ExecStartPost=-/usr/bin/systemctl --user enable --now iva-update-check.timer
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
