# deploy/iva-bitrix-sync.service
[Unit]
Description=Iva Bitrix24 read-only task synchronization
After=network-online.target iva.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__NODE_BIN__ --import=tsx --env-file=.env scripts/bitrix-sync.ts --daily
NoNewPrivileges=true
PrivateTmp=true
