# deploy/iva-memory-doctor.service
[Unit]
Description=Iva memory doctor (vault health + git push)
After=network-online.target iva.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__NODE_BIN__ --env-file=.env scripts/memory/doctor.ts
